# backups/omarchy-dotfiles-20260214-111647/.config/systemd/user/voxtype.service.d/backend.conf
[Service]
ExecStart=
ExecStart=/usr/lib/voxtype/voxtype-vulkan daemon
